# templates/etc/profile.d/subgit.sh.j2
export PATH=$PATH:{{subgit_home}}/subgit/bin
